# deploy/systemd/vps-guardian-compose-backup-freshness.timer
[Unit]
Description=Monitor VPS Guardian Compose backup freshness

[Timer]
OnBootSec=30min
OnUnitActiveSec=15min
Persistent=true
Unit=vps-guardian-compose-backup-freshness.service

[Install]
WantedBy=timers.target
